# crontab
0 * * * * /root/digitalocean_availability.sh
